# software/kernel/tpr.service
[Unit]
Description=Tpr Device Manager
Requires=multi-user.target
After=multi-user.target

[Service]
Type=simple
ExecStart=/etc/systemd/system/tpr.d/tpr.init.d start
ExecStartPost=
ExecStop=/etc/systemd/system/tpr.d/tpr.init.d stop
KillMode=process
IgnoreSIGPIPE=no
RemainAfterExit=yes
StandardOutput=syslog
StandardError=inherit

[Install]
WantedBy=multi-user.target
